# src/root/usr/lib/systemd/system/brutalflea.service
[Unit]
Description=Hello World Service

[Service]
ExecStart=/usr/local/sbin/brutalflea.sh

[Install]
WantedBy=multi-user.target